# mosaic/rootfs/etc/services.d/mosaic/finish
#!/usr/bin/execlineb -S1
# shellcheck shell=bash
# ==============================================================================
# Take down the S6 supervision tree when Mosaic fails
# ==============================================================================

/run/s6/basedir/bin/halt
